# Arch Linux/systemd/optimize-pacman.service
[Unit]
Description=Optimize pacman database

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman-optimize
