# Linux/systemd_user/end-windows7vm-maint.service
[Unit]
Description=End Windows 7 VM 5 AM on Monday

[Service]
Type=oneshot
ExecStart=/usr/bin/vboxmanage controlvm Windows7 acpipowerbutton
