# Linux/systemd_user/start-windows7vm-maint.service
[Unit]
Description=Start Windows 7 VM 3 AM on Monday

[Service]
Type=oneshot
ExecStart=/usr/bin/vboxheadless -s Windows7
