# Linux/systemd/reboot-windows-10.service
[Unit]
Description=Reboot to Windows 10

[Service]
Type=oneshot
ExecStart=/usr/bin/grub-reboot 2
ExecStartPost=/usr/bin/reboot
